# addons/hadocs/rootfs/etc/services.d/hadocs/run
#!/usr/bin/with-contenv bashio

set -e

PROJECT_NAME="$(bashio::config 'project_name')"
OUTPUT_DIRECTORY="$(bashio::config 'output_directory')"

mkdir -p "${OUTPUT_DIRECTORY}"

export HADOCS_HA_URL="http://supervisor/core"
export HADOCS_TOKEN="${SUPERVISOR_TOKEN}"
export HADOCS_PROJECT_NAME="${PROJECT_NAME}"
export HADOCS_OUTPUT_DIR="${OUTPUT_DIRECTORY}"
export HADOCS_CACHE_DIR="/tmp/hadocs-cache"
export HADOCS_CONFIG_FILE="/config/hadocs.json"

bashio::log.info "Starting HADocs scan"
bashio::log.info "Output directory: ${OUTPUT_DIRECTORY}"

python -m src.hadocs.cli.main generate

bashio::log.info "HADocs scan completed"